# util/systemd/rustic-forget@.service
[Unit]
Description=rustic --use-profile %I forget

[Service]
KillSignal=SIGINT
ExecStart=/usr/bin/rustic --use-profile %I forget
